# frontend/deploy/systemd/contexthelix-console-backend.service
[Unit]
Description=ContextHelix Console Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/data/Autopoiesis-AgentSys/frontend
Environment=PYTHONUNBUFFERED=1
Environment=CONTEXTHELIX_ARTIFACT_ROOT=/data/netops-runtime
Environment=CONTEXTHELIX_REPO_ROOT=/data/Autopoiesis-AgentSys
Environment=CONTEXTHELIX_FRONTEND_DIST=/data/Autopoiesis-AgentSys/frontend/dist
ExecStart=/data/Autopoiesis-AgentSys/frontend/.venv/bin/uvicorn gateway.app.main:app --app-dir /data/Autopoiesis-AgentSys/frontend --host 127.0.0.1 --port 8026
Restart=always
RestartSec=3
TimeoutStopSec=8
KillMode=mixed

[Install]
WantedBy=multi-user.target
